# fan.service
[Unit]
Description=PWM fan controller
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/alex/raspi-tweaks
ExecStart=/usr/bin/python3 /home/alex/raspi-tweaks/fan_controller.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target